# location-server.service.example
[Unit]
Description=Feishu Task Location Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/guijie/feishu-task-listener
ExecStart=/usr/bin/python3.11 /home/guijie/feishu-task-listener/location_server.py --host 0.0.0.0 --port 8000 --lark-cli lark-cli
Restart=always
RestartSec=5
User=guijie
Environment=PYTHONUNBUFFERED=1
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/home/guijie/.local/bin

[Install]
WantedBy=multi-user.target
